# roles/poll/files/poll.service
[Unit]
Description=Poll service file

[Service]
User=root
Environment=REDIS_HOST
ExecStart=python3 /home/user/poll/app.py
Restart=always